# profile.d/cargo.sh
# cargo.sh --- configure the rust cargo setup

export RUST_SRC_PATH=$HOME/src/rust/src
export PATH=$PATH:$HOME/.cargo/bin
