# profile.d/fasd.bash
#!/bin/bash
# fasd.bash --- Initialize the fasd command line navigator
# https://github.com/clvv/fasd

# Mike Barker <mike@thebarkers.com>
# March 5th, 2019

eval "$(fasd --init auto)"

