# profile.d/history.bash
# History configuration
export HISTSIZE=1000
export HISTIGNORE="&b:[bf]g:exit"
#export HISTIGNORE="&:bg:fg:ll:h"
export HISTTIMEFORMAT="$(echo ${BCyan})[%d/%m %H:%M:%S]$(echo ${NC}) "
export HISTCONTROL=ignoredups
