# profile.d/history.bash
# History configuration
export HISTSIZE=1000
export HISTIGNORE="&:?:??:exit:pwd:clear"
export HISTTIMEFORMAT="[%d/%m %H:%M:%S] "
export HISTCONTROL=ignoreboth:erasedups

alias h=history
