# profile.d/history.bash
# History configuration
export HISTSIZE=1000
export HISTIGNORE="&b:[bf]g:exit"
