# profile.d/hosts.sh
if [ -f $HOME/.hosts ]; then
    export HOSTFILE=$HOME/.hosts    # Put a list of remote hosts in ~/.hosts
fi
