# profile.d/hunspell.bash
#!/bin/bash
# spelling setup -- on osx, use brew install hunspell

export DICTIONARY=en_US
