# profile.d/options.sh
# Manipulate shell options
set -o notify	   # report status of terminated background jobs immediately
set -o noclobber   # prevent redirection from overwriting files
set -o ignoreeof   # do not exit on eof
