# profile.d/rvm.bash
# rvm setup
RVM_SCRIPT="${HOME}/.rvm/scripts/rvm"
[ -f "${RVM_SCRIPT}" ] && source "${RVM_SCRIPT}"
