# profile.d/virtualenv.bash
# virtualenv setup
if [[ -x /usr/local/bin/virtualenv ]]; then
    if [[ -x /usr/local/bin/virtualenvwrapper.sh ]]; then
	. /usr/local/bin/virtualenvwrapper.sh
    fi
fi
